# wallpapers/wallpaper.timer
[Unit]
Description=Wallpaper Change Timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=2h
Unit=wallpaper.service

[Install]
WantedBy=timers.target

